# .agents/idea-router/wake.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
Example launchd plist for FSEvents-driven router wake on macOS.

Copy to ~/Library/LaunchAgents/com.sirsi.idea-router.plist, edit the
absolute paths and PATH entries for your machine, then:

  launchctl load ~/Library/LaunchAgents/com.sirsi.idea-router.plist

Behavior:
  - launchd watches state.json, items/, and proposals/ for ANY change
  - On change, spawns ONE dispatch pass (`sirsi router run --once`)
  - Process exits; launchd returns to watching
  - ThrottleInterval=10 prevents refire within 10s (avoids loops when
    a dispatch writes back to state.json and re-triggers the watch)
  - ZERO idle process — no polling, no daemon, no cron

This replaces the prior push-model daemon (polling every 1s, ~86,400
reads/day, ~24h CPU drain to dispatch zero useful work because the
target binary was missing). The FSEvents version dispatches in
milliseconds on file create/modify and burns zero cycles when idle.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.sirsi.idea-router</string>
  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/sirsi</string>
    <string>router</string><string>run</string>
    <string>--once</string>
    <string>--target</string><string>all</string>
  </array>
  <key>WorkingDirectory</key><string>/ABSOLUTE/PATH/TO/sirsi-pantheon</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>SIRSI_ROUTER_NOTIFY</key><string>1</string>
    <key>PATH</key><string>/Users/YOU/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>WatchPaths</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/.agents/idea-router/state.json</string>
    <string>/ABSOLUTE/PATH/TO/.agents/idea-router/items</string>
    <string>/ABSOLUTE/PATH/TO/.agents/idea-router/proposals</string>
  </array>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/ABSOLUTE/PATH/TO/.agents/idea-router/logs/launchd.out.log</string>
  <key>StandardErrorPath</key><string>/ABSOLUTE/PATH/TO/.agents/idea-router/logs/launchd.err.log</string>
</dict>
</plist>
